# hive/playbooks/roles/base/files/NetworkManager-wait-online.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/nm-online -q --timeout=30
SuccessExitStatus=1